# usr/share/bgscripts/examples/shares-keepalive.cron
*/3	*	*	*	*	root	/usr/share/bgscripts/shares.sh --all --keepalive 1>/dev/null 2>&1
