# usr/share/bgscripts/build/debian-bgscripts-core/postinst
# deb core postinst 2017-07-21
{
if test "$1" = "configure";
then
   if test "$( ps --no-headers -o comm 1 )" = "systemd";
   then
      install -m 0644 -o root -p -t "/lib/systemd/system/" "/usr/share/bgscripts/inc/systemd/dnskeepalive.service" || :
      install -m 0644 -o root -p -t "/lib/systemd/system-preset/" "/usr/share/bgscripts/inc/systemd/80-dnskeepalive.preset" || :
      systemctl daemon-reload &
   fi
   systemctl --no-reload preset dnskeepalive.service || :

   # install symlink for doc
   test ! -h /usr/share/bgscripts/doc && ln -s ../doc/bgscripts /usr/share/bgscripts/doc
fi
} 1>/dev/null 2>&1
