# usr/share/bgscripts/build/debian-bgscripts-core/prerm
# deb core prerm 2017-06-08
# Reference: https://www.debian.org/doc/debian-policy/ch-maintainerscripts.html
{
if test "$1" = "remove";
then
   # Package  removal, not upgrade
   systemctl --no-reload disable --now dnskeepalive.service || :
   rm -f /lib/systemd/system/dnskeepalive.service || :
   rm -f /lib/systemd/system-preset/80-dnskeepalive.preset || :
fi
} 1>/dev/null 2>&1
