# usr/share/bgscripts/build/debian-bgscripts-core/prerm
# deb core prerm 2018-02-23
# Reference: https://www.debian.org/doc/debian-policy/ch-maintainerscripts.html
{
if test "$1" = "remove" ;
then
   # Total uninstall

   # remove systemd files
   systemctl --no-reload disable --now dnskeepalive.service || :
   rm -f /lib/systemd/system/dnskeepalive.service || :
   rm -f /lib/systemd/system-preset/80-dnskeepalive.preset || :

   # remove doc symlink, only for deb
   test -L /usr/share/bgscripts/doc && unlink /usr/share/bgscripts/doc

fi
} 1>/dev/null 2>&1

exit 0
