# usr/share/bgscripts/build/debian-bgscripts/prerm
# deb prerm 2017-09-16
{
if test "$1" = "remove";
then
   # total uninstall

   # Remove mimetype definitions
   for user in root ${SUDO_USER} Bgirton bgirton bgirton-local;
   do
      which xdg-mime && {
         su "${user}" -c "xdg-mime uninstall /usr/share/bgscripts/gui/x-rdp.xml &"
      }
   done

   # remove systemd files
   systemctl --no-reload disable --now monitor-resize.service || :
   rm -f /lib/systemd/system/monitor-resize.service || :
   rm -f /lib/systemd/system-preset/80-monitor-resize.preset || :

fi
} 1>/dev/null 2>&1
exit 0
